# crontab
* * * * * echo "Crontab is working - watchdog 1" >> /tmp/cronjoblog
